# scripts/mac-mini/com.devbox.ddns.plist
<?xml version="1.0" encoding="UTF-8"?>
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.devbox.ddns</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/sccarey/devbox/ddns.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>300</integer>
  <key>StandardOutPath</key>
  <string>/Users/sccarey/Library/Logs/ddns.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/sccarey/Library/Logs/ddns.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/sccarey</string>
  </dict>
</dict>
</plist>
