# services/whoeverwants-dev.service
[Unit]
Description=WhoeverWants Development Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/whoeverwants
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=whoeverwants-dev
Environment="NODE_ENV=development"
Environment="PATH=/usr/bin:/usr/local/bin:/home/ubuntu/.nvm/versions/node/v20.11.1/bin"

[Install]
WantedBy=multi-user.target